# edutrack-lifestyle.service
[Unit]
Description=EduTrack Lifestyle Suggestion API (Flask)
After=network.target

[Service]
User=www-data
WorkingDirectory=/var/www/html/edutrack
ExecStart=/var/www/html/edutrack/venv/bin/gunicorn --workers 2 --bind 127.0.0.1:5003 lifestyle_api:app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
